# systemd/ecoflow-ble-nutd.service
[Unit]
Description=EcoFlow BLE NUT-compatible daemon
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ecoflow-ble-nutd -config /etc/ecoflow-ble-nutd.conf
Restart=on-failure
RestartSec=5
DynamicUser=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run
MemoryMax=64M

[Install]
WantedBy=multi-user.target
